# AUTOMATIONS/agent/autonomy/schedules/com.claude.schedule.auto_local_biz_openclaw_nationwide_9569.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.auto_local_biz_openclaw_nationwide_9569</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt" &amp;&amp; claude -p "You are the LOCAL BIZ autonomy agent for PRINTMAXX venture '\''OpenClaw Nationwide'\''. Working directory: /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt  Your job runs every 4h. Each cycle: 1. DISCOVER: Search for local businesses in rotating cities that have no website or terrible websites. Use DuckDuckGo search, Yelp scraping, Google Maps data. 2. GRADE: Score each business website A-F based on mobile responsiveness, SSL, load time, design quality. 3. BUILD PREVIEW: For F and D grade businesses, auto-generate a modern preview website. 4. DEPLOY: Deploy preview to surge.sh with unique URL. 5. OUTREACH: Generate cold email with preview link. Queue for sending. 6. TRACK: Check outreach responses, follow up on warm leads.  City rotation: cycle through US cities automatically. Can also target international markets. Rules: All files stay in /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt  AGENTIC EXECUTION: Do NOT stop after one step. Keep going until the task is COMPLETE or you hit a blocker. - Read state -> Plan -> Execute -> Verify -> If more work needed, CONTINUE - If you find sub-tasks while working, do them now, don'\''t defer - Only stop when: (a) task fully complete, (b) human action required, (c) 15+ tool calls made (safety limit) - Write completion status to AUTOMATIONS/agent/autonomy/auto_local_biz_openclaw_nationwide_9569/output/report_20260328.md" --dangerously-skip-permissions --model claude-sonnet-4-6 >> "/Users/macbookpro/.claude/logs/auto_local_biz_openclaw_nationwide_9569.log" 2>&amp;1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/macbookpro/.local/bin</string>

    </dict>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>StandardOutPath</key>
    <string>/Users/macbookpro/.claude/logs/auto_local_biz_openclaw_nationwide_9569.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/.claude/logs/auto_local_biz_openclaw_nationwide_9569.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>