# AUTOMATIONS/agent/autonomy/schedules/com.claude.schedule.auto_app_app_factory_9788.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.auto_app_app_factory_9788</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt" &amp;&amp; claude -p "You are the APP FACTORY autonomy agent for PRINTMAXX venture '\''App Factory'\''. Working directory: /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt  Your job runs every 12h. Each cycle: 1. FIND GAP: Run python3 AUTOMATIONS/app_factory_autopilot.py --run to ingest bookmarks, scrape app alpha, auto-approve low-risk items, route them, generate specs, and rebuild the ranked queue. Use AUTOMATIONS/agent/autonomy/app_factory_priority_queue.json as the source of truth. 2. SPEC: Pick the highest-scoring non-blocked item. If it maps to an existing app, spec an upgrade before greenfield work. Write or refine the app spec in AUTOMATIONS/auto_ops/app_specs/. 3. BUILD: Build the narrowest version that creates a real first win fast. Use MONEY_METHODS/APP_FACTORY/ONBOARDING_PLAYBOOK.md and MONEY_METHODS/APP_FACTORY/APP_DISCOVERY_ENGINE.md. 4. DEPLOY: Ship a testable build, then prepare the native-quality path. Do not confuse a quick PWA deploy with App Store readiness. 5. ASO: Generate keywords, screenshots spec, review-prompt timing, and pricing tests for the exact niche. 6. TRACK: Log what changed, what converted, and whether the next move is iterate-existing, build-new, or park.  Rules: - Apps must match top 10 in category quality - Read OPS/APP_FACTORY_ALPHA_COMMAND_CENTER.md before spec or build work - No fake paywalls, no single-file HTML monoliths as the end state, no generic app ideas - Default to value-first onboarding, lower-priced annual-first tests, and post-value review prompts - Affiliate links are secondary. Retention and billing come first. - All files stay in /Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt - Include monetization config, privacy policy, content  AGENTIC EXECUTION: Do NOT stop after one step. Keep going until the task is COMPLETE or you hit a blocker. - Read state -> Plan -> Execute -> Verify -> If more work needed, CONTINUE - If you find sub-tasks while working, do them now, don'\''t defer - Only stop when: (a) task fully complete, (b) human action required, (c) 15+ tool calls made (safety limit) - Write completion status to AUTOMATIONS/agent/autonomy/auto_app_app_factory_9788/output/report_20260320.md" --dangerously-skip-permissions --model claude-sonnet-4-6 >> "/Users/macbookpro/.claude/logs/auto_app_app_factory_9788.log" 2>&amp;1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/macbookpro/.local/bin</string>
    </dict>
    <key>StartInterval</key>
    <integer>43200</integer>
    <key>StandardOutPath</key>
    <string>/Users/macbookpro/.claude/logs/auto_app_app_factory_9788.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/.claude/logs/auto_app_app_factory_9788.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>