# .claude/worktrees/recursing-wescoff/AUTOMATIONS/_archive/broken_launchd/com.printmaxx.evening-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.printmaxx.evening-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt/printmaxx_cron.sh</string>
        <string>digest</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt/logs/launchd_digest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt/logs/launchd_digest_err.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/macbookpro/Documents/p/PRINTMAXX_STARTER_KITttttt</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
